# element/standalone/skeleton/config/etc/supervisor.conf
[program:www.mywebsite.com]
command=/path/to/virtualenv/bin/gunicorn -w 3 -b 127.0.0.1:5000 wsgi:app
autostart=true
autorestart=true
startsecs=2
user=nobody
directory=/path/to/virtualenv/site/src